# system/unit_files/rotation.timer
# /etc/systemd/system/rotation.timer
[Unit]
Description=rotation timer

[Timer]
OnStartupSec=20

[Install]
WantedBy=basic.target
